# .docker/s6-overlay/s6-rc.d/svc-scheduler/run
#!/command/with-contenv sh
cd /app
exec php artisan schedule:work --no-interaction
